# pi/systemd/dnsmasq.service.d/override.conf
[Unit]
# dnsmasq darf erst starten, wenn wlan0 als AP läuft und 192.168.4.1 hat.
After=hostapd.service network-online.target sys-subsystem-net-devices-wlan0.device
Wants=hostapd.service sys-subsystem-net-devices-wlan0.device
Requisite=sys-subsystem-net-devices-wlan0.device

[Service]
# Falls wlan0 beim ersten Versuch noch nicht bereit ist: warten & neu versuchen.
ExecStartPre=/bin/sh -c 'for i in $(seq 1 20); do ip -4 addr show wlan0 | grep -q "192.168.4.1" && exit 0; sleep 1; done; exit 1'
Restart=on-failure
RestartSec=3
