# pi/systemd/video-token-ap.service
[Unit]
Description=Guest_Video_Token AP-Modus beim Boot
After=network.target dhcpcd.service video-token-bootmode.service
Wants=dhcpcd.service
ConditionPathExists=!/run/video-token-client-mode

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/switch-mode ap
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
