# pi/systemd/video-token-admin.service
[Unit]
Description=Guest_Video_Token Admin API
After=network.target
ConditionPathExists=!/run/video-token-client-mode

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/local/sbin/admin-server.py
Restart=on-failure
RestartSec=3
User=root

[Install]
WantedBy=multi-user.target
